# etc/supervisord.d/opensm.ini
[program:opensm]
command=/sbin/opensm
redirect_stderr=true
#stdout_logfile = /var/log/supervisor/%(program_name)s.log
stdout_syslog=true

